# scripts/deploy/amps-startup.service
[Unit]
Description=Amps-app startup service

[Service]
Type=exec
ExecStart=amps start

[Install]
WantedBy=multi-user.target