# scripts/deploy/amps-startup.service
[Unit]
Description=Amps-app startup service

[Service]
Type=exec
ExecStart=/usr/local/bin/amps-startup.sh

[Install]
WantedBy=multi-user.target